# docker/root/etc/services.d/fightarr/run
#!/usr/bin/with-contenv bash

cd /app
exec s6-setuidgid abc dotnet Fightarr.Web.dll
